# bakery-ops/com.hotcrush.bakery-ops.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hotcrush.bakery-ops</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/weiliangshao/hot/bakery-ops/start-daemon.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/Users/weiliangshao/hot/bakery-ops/logs/daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/weiliangshao/hot/bakery-ops/logs/daemon.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/weiliangshao/hot/bakery-ops</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
